# docs/deploy/mermaduckle.service
[Unit]
Description=Mermaduckle server (example systemd unit)
After=network.target

[Service]
# Create a dedicated `mermaduckle` user and place the binary in /opt/mermaduckle
User=mermaduckle
Group=mermaduckle
WorkingDirectory=/opt/mermaduckle
ExecStart=/opt/mermaduckle/mermaduckle-server
Environment=DATABASE_PATH=/var/lib/mermaduckle/app.db
Environment=PORT=3001
Restart=on-failure
RestartSec=5
LimitNOFILE=65536
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

Notes:
- Build the release binary on the host or cross-compile, then copy to `/opt/mermaduckle`.
- Adjust `User` and paths for your environment. Enable and start with `systemctl enable --now mermaduckle`.
